# playbooks/roles/website/files/supervisord/app.ini
[program:ilmnuri]
command=uwsgi --master --emperor /usr/share/nginx/html --die-on-term
directory=/usr/share/nginx/html/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/app-evolvecyber.log
redirect_stderr=true
stopsignal=QUIT
